# cron.d/apt-mirror
 0 4 * * * apt-mirror /run_apt-mirror.sh
